# packages/piratepython/piratepython.service
[Unit]
Description=PiratePython main.py
After=network.target
Wants=network.target

[Service]
ExecStart=/usr/bin/python3 /mtp/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
