# packages/piratepython/pp-restart.service
[Unit]
Description=PiratePython Restarter

[Service]
Type=oneshot
ExecStart=/bin/systemctl restart piratepython

[Install]
WantedBy=multi-user.target
