# packages/piratepython/pp-restart.path
[Path]
PathChanged=/mtp/main.py

[Install]
WantedBy=multi-user.target
