# packages/qmldemo/qmldemo.service
[Unit]
Description=QML Demo Scene
After=network.target
Wants=network.target

[Service]
Environment=QT_QPA_EVDEV_KEYBOARD_PARAMETERS=grab=1
ExecStart=/opt/qt/bin/qmlscene /opt/qt/bin/qmldemo.qml
Restart=on-failure

[Install]
WantedBy=multi-user.target
