# packages/pymtpd/pymtpd.service
[Unit]
Description=MTP Responder
After=network.target
Wants=network.target

[Service]
Environment=PYTHONPATH=/opt/pymtpd/lib/python3.6/site-packages
ExecStart=/opt/pymtpd/bin/pymtpd -n "Pirate Python" --udc 20980000.usb -s storage /mtp -P "Pirate Python"
Restart=on-failure

[Install]
WantedBy=multi-user.target
