# packages/qmlrss/qmlrss.service
[Unit]
Description=QML RSS Reader

[Service]
Environment=QT_QPA_EVDEV_KEYBOARD_PARAMETERS=grab=1
ExecStart=/opt/qt/bin/qmlrss https://feeds.bbci.co.uk/news/technology/rss.xml
Restart=on-failure

[Install]
WantedBy=multi-user.target
